# data/packages/mesk-units/mesk-units.install
post_install()
{
systemctl enable local.d
systemctl enable rc-local
}
