# data/packages/mesk-units/rc-local.service
[Unit]
Description=Run rc.local
Before=default.target
[Service]
ExecStart=/etc/rc.local
Type=oneshot

[Install]
RequiredBy=default.target

