# data/overlay/etc/local.d/linux_logo.start
#!/bin/bash
mv /etc/issue /etc/issue.tmp
linux_logo -L classic-simp > /etc/issue
cat /etc/issue.tmp >> /etc/issue
rm /etc/issue.tmp
