# data/overlay/etc/local.d/mesk_user.start
#!/bin/bash

sleep 1 && chown -hR mesk:users /home/mesk/ &> /dev/null
