# data/overlay/etc/local.d/sound.start
if ! [[ -e /proc/asound ]]
then
/etc/init.d/oss start &
fi
